# apps/openclaw-web-platform/deploy/systemd/openclaw-local-console.service
[Unit]
Description=SoloCore Console Cloud Broker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/openclaw/mission-control
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=OPENCLAW_DISABLE_AGENT_AUTOSTART=1
Environment=NODE_OPTIONS=--max-old-space-size=384
ExecStart=/usr/bin/node /opt/openclaw/mission-control/dist-server/server.js
Restart=always
RestartSec=5
MemoryAccounting=true
MemoryMax=650M
TasksMax=256
User=openclaw
Group=openclaw

[Install]
WantedBy=multi-user.target
